# agent/replicator_svc.conf
description "Cargo Lazy Replicator"
author "shnadgow@in.ibm.com"

start on runlevel [2345]
stop on runlevel [0156]

script
CMD
end script

